# packaging/debian/vantagenet-daemon.service
[Unit]
Description=VantageNet Network Discovery Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/vantagenet-daemon
Restart=always
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal

# Security hardening - allow only required capabilities
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN

# Security restrictions
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/etc/vantagenet

[Install]
WantedBy=multi-user.target
